# cookbooks/main/templates/default/rails_app.unicorn_init.conf.erb
start on startup
stop on shutdown

chdir <%= node[:rails_app][:www_app_path] %>/current
env RAILS_ENV=production
normal exit 0 QUIT USR2

pre-start script
  su -c "bundle exec unicorn -c config/unicorn.rb -D" <%= node[:user][:name] %>
end script

post-stop script
  exec kill `cat <%= node[:rails_app][:unicorn_pid_path] %>`
end script
